# deploy/caddy.service
[Unit]
Description=Caddy fronting Clove on :8443 (cert from acme.sh / DNSPod DNS-01)
Documentation=https://caddyserver.com/docs/
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
User=winbeau
Group=winbeau
# /usr/local/bin/caddy is the binary you copied from /tmp/caddy
ExecStart=/usr/local/bin/caddy run --environ --config /etc/caddy/Caddyfile
ExecReload=/usr/local/bin/caddy reload --config /etc/caddy/Caddyfile --force
TimeoutStopSec=5s
LimitNOFILE=1048576
PrivateDevices=true
ProtectHome=read-only
ProtectSystem=full
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
