# deploy/clove.service
[Unit]
Description=Clove Claude reverse proxy
Documentation=https://github.com/mirrorange/clove
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=winbeau
Group=winbeau
WorkingDirectory=/home/winbeau/clove
EnvironmentFile=/home/winbeau/clove/.env

# Run via the uv-managed venv directly. We deliberately don't go through
# `uv run` here because uv writes a lockfile at $UV_CACHE_DIR (~/.cache/uv
# by default) on every invocation, which collides with ProtectHome=read-only.
# uv is only needed at sync time (`make install-dev`); at run time we just
# invoke the venv's interpreter.
ExecStart=/home/winbeau/clove/.venv/bin/python -m app.main

Restart=on-failure
RestartSec=5

# Hardening — relaxed enough to let DATA_FOLDER + log file work.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/winbeau/clove /home/winbeau/.clove
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
